# support/systemd/prefault.service
[Unit]
Description=prefault files into memory

[Service]
Type=exec
ExecStart=/usr/bin/prefault mlock

[Install]
WantedBy=multi-user.target